# WebApp/PIMS/settings/adminapp.socket
[Unit]
Description=gunicorn socket to trigger adminapp.service

[Socket]
ListenStream=/run/adminapp.sock

[Install]
WantedBy=sockets.target
